# crontab
21 */2 * * * ~/.autosave/autosave.sh &> /tmp/autosave.log
